# launchd/com.tradebot.aapl-10am-once.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.aapl-10am-once</string>

    <key>WorkingDirectory</key>
    <string>/Users/zsayed/personal/trading-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/zsayed/personal/trading-bot</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zsayed/personal/trading-bot/scripts/run-aapl-10am-once.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Day</key>
        <integer>5</integer>
        <key>Month</key>
        <integer>6</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zsayed/.tradebot/aapl-10am-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zsayed/.tradebot/aapl-10am-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
