# launchd/com.tradebot.market-open.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.market-open</string>

    <key>WorkingDirectory</key>
    <string>/Users/zsayed/personal/trading-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/zsayed/personal/trading-bot</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zsayed/personal/trading-bot/.venv/bin/python</string>
        <string>-m</string>
        <string>tradebot.daemon</string>
        <string>--no-server</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
        <key>Day</key>
        <integer>5</integer>
        <key>Month</key>
        <integer>6</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zsayed/.tradebot/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zsayed/.tradebot/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
